# scripts/rc.local
#!/bin/bash
bash /boot/setup-cowrie.sh
exit 0
